# 2.projects/sglang-inference/src/sglang-worker.service
[Unit]
Description=SGLang Worker Service
After=network.target cloud-final.service
Wants=cloud-final.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/sglang
Environment="TORCHINDUCTOR_CACHE_DIR=/opt/sglang/torch_compile_cache/"
Environment="PYTHONPATH=/opt/sglang/source:$PYTHONPATH"
ExecStart=/opt/sglang/venv/bin/python /opt/app/run_worker.py --gpu-id 0 --model MODEL_PATH --router-url http://10.0.0.100:8000 EXTRA_ARGS
Restart=always
RestartSec=10
StandardOutput=append:/opt/sglang/logs/sglang.log
StandardError=append:/opt/sglang/logs/sglang.log

[Install]
WantedBy=multi-user.target